# autoupdater/crontab
SHELL=/bin/sh

* * * * * curl http://web:8080/autoupdate_match_result