# autoupdater/crontab
SHELL=/bin/sh

* * * * * curl http://web/autoupdate_match_result