# autoupdater/crontab
SHELL=/bin/sh

* * * * * curl -s http://web:8080/autoupdate_match_result 2>&1
