# backup-daemon/crontab
SHELL=/bin/sh

0 * * * * /app/backup.sh 2>&1
